# roles/borg/templates/borgmatic.timer.j2
[Unit]
Description=Run borgmatic backup

[Timer]
OnCalendar=daily
RandomizedDelaySec=2h

[Install]
WantedBy=timers.target
