# roles/borg/templates/borgmatic.service.j2
[Unit]
Description=borgmatic backup

[Service]
Type=oneshot
ExecStart={{ borgmatic_path }} -v 1
